# deploy/sso.service
[Unit]
Description=Dalang SSO — self-hosted OAuth 2.0 / OIDC provider
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/dalang-sso
EnvironmentFile=/opt/dalang-sso/.env
ExecStart=/opt/dalang-sso/sso
Restart=on-failure
RestartSec=2
# Bind address comes from .env (SSO_BIND_ADDR). Behind Cloudflare "Flexible"
# SSL the app serves plain HTTP directly on 0.0.0.0:80; behind a local TLS
# proxy set it to 127.0.0.1:8080 instead.

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dalang-sso/data
PrivateTmp=true
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
